# supervisord.conf
[supervisord]
nodaemon=true

[program:celery_worker]
command=celery -A celery_app worker --loglevel=INFO
autorestart=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0

[program:fastapi]
command=uvicorn main:app --host %(ENV_HOST)s --port %(ENV_PORT)s --use-colors
autostart=true
autorestart=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0

[program:flower]
command=celery -A celery_app flower --port=%(ENV_FLOWER_PORT)s
autostart=true
autorestart=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0